# printbrigata.service
[Unit]
Description=Printbrigata Rechnungstool
After=network.target

[Service]
User=root
# The user is running in /root/rechnungs-tool
WorkingDirectory=/root/rechnungs-tool
Environment="PATH=/root/rechnungs-tool/venv/bin:/usr/local/bin:/usr/bin:/bin"

# command to start the app
# Note: Using 'venv' (not .venv) because the user created it manually as 'venv'
ExecStart=/root/rechnungs-tool/venv/bin/python -m streamlit run Preislisten.py --server.port 8501 --server.address 127.0.0.1 --server.headless true

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
